# com.ai-file-renamer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-file-renamer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jn/Apps/ai-file-renamer/venv/bin/python3</string>
        <string>/Users/jn/Apps/ai-file-renamer/file_renamer.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/jn/Apps/ai-file-renamer/venv/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/jn/Apps/ai-file-renamer</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/jn/Apps/ai-file-renamer</string>
    <key>StandardErrorPath</key>
    <string>/Users/jn/Library/Logs/ai-file-renamer/error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/jn/Library/Logs/ai-file-renamer/output.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist> 